# Assignment_04/99-keyboard.rules
ACTION=="add", \
SUBSYSTEM=="input", \
ENV{ID_BUS}="usb", \
ENV{ID_INPUT_KEYBOARD}="1", \
RUN+="/home/blfs/Little_Pinguin/Assignment_04/insert.sh"
